# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Enable and start the service on fresh install
        if [ -z "$2" ]; then
            deb-systemd-helper enable cdi-mp.service >/dev/null 2>&1 || true
        fi
        deb-systemd-helper start cdi-mp.service >/dev/null 2>&1 || true
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        exit 0
        ;;
    *)
        exit 0
        ;;
esac

exit 0
